# deploy/danci-ota-service.service
[Unit]
Description=Danci OTA Update Service
Documentation=https://github.com/heartcoolman/vocabulary-learning-app
After=network.target docker.service
Requires=docker.service

[Service]
Type=simple
ExecStart=/opt/danci/scripts/prod-ota-service.sh
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
